# config/crontab
*/2 * * * * root echo Scheduled Generation && /Deej-AI-Auto/config/generate.sh
